# docker/etc/s6-overlay/s6-rc.d/listener/run
#!/usr/bin/env sh
with-contenv
cd /var/www/tf2_live_stats; bundle exec /var/www/tf2_live_stats/script/listener
