# docker/etc/s6-overlay/s6-rc.d/thin/run
#!/usr/bin/env sh
with-contenv
cd /var/www/tf2_live_stats; bundle exec thin -C /var/www/tf2_live_stats/config/thin.yml start